# crontab
* * * * * /venv/bin/planningpod-stats >> /logfile
