# crontab
# Run one minute after every hour
1 * * * * /venv/bin/planningpod-stats >> /logfile
